# consul_agent/templates/etc/systemd/system/consul-agent-backup.service.j2
{{ ansible_managed | comment }}

[Unit]
Description=Script for regular consul state backups

[Service]
ExecStart={{ _consul_agent_backup_script_path }}
SyslogIdentifier=%p
Type=oneshot
